# node/supervisord.conf
[program:twitterpoll]
command = node /home/dotcloud/current/server.js
