# node/supervisord.conf
[program:twitterpoll]
command = node /home/dotcloud/current/twitter-poll.js
